# wifi-browser.service
[Unit]
Description=Wifi Browser

[Service]
ExecStart=/usr/bin/node /home/root/wifi-browser-master/app.js

[Install]
WantedBy=multi-user.target

